# docker/_shared/ai-lab-secrets@.service
# Instanced unit: ai-lab-secrets@chat or ai-lab-secrets@gateway
# Fetches lab/* secrets into a tmpfs .env BEFORE the compose stack starts.
# Install: copy to /etc/systemd/system/ ; edit REPO_DIR/region if needed.
[Unit]
Description=AI Lab secrets bootstrap (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=AI_LAB_ROLE=%i
Environment=AWS_DEFAULT_REGION=us-east-1
Environment=AI_LAB_COMPOSE=/opt/ai-lab/repo/docker/%i-host
ExecStart=/opt/ai-lab/repo/docker/_shared/secrets-bootstrap.sh

[Install]
WantedBy=multi-user.target
