# docker/_shared/ai-lab-stack@.service
# Instanced unit: ai-lab-stack@chat or ai-lab-stack@gateway
# Brings up the docker compose stack for this host role, after secrets are seeded.
# Install: copy to /etc/systemd/system/ ; `systemctl enable --now ai-lab-stack@gateway`
[Unit]
Description=AI Lab docker compose stack (%i)
After=docker.service ai-lab-secrets@%i.service
Requires=docker.service ai-lab-secrets@%i.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/ai-lab/repo/docker/%i-host
# The compose BUILD (buildx) is a separate client from dockerd; it needs the
# Squid proxy in ITS env to pull base-image manifests, or it dials docker.io
# directly and the app SG drops it. The daemon proxy (docker.service.d) only
# covers `docker pull`, not `compose --build`. Load the host proxy env here.
# `-` = tolerant if absent (e.g. networkfirewall egress mode).
EnvironmentFile=-/etc/environment
ExecStart=/usr/bin/docker compose up -d --build
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
